# invit-watcher.path
[Path]
PathModified=/root/go/src/github.com/invitation/bin/

[Install]
WantedBy=multi-user.target
